# debian/envoyctl-apply.service
[Service]
Type=oneshot
EnvironmentFile=-/etc/default/envoyctl
WorkingDirectory=%E{ENVOYCTL_WORKDIR}

ExecStart=/usr/bin/envoyctl apply \
  --config-dir %E{ENVOYCTL_CONFIG_DIR} \
  --out-dir %E{ENVOYCTL_OUT_DIR} \
  --install-path %E{ENVOYCTL_INSTALL_PATH}

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Allow writes only where needed
ReadWritePaths=%E{ENVOYCTL_WORKDIR} %E{ENVOYCTL_INSTALL_DIR}

# Networking is needed only if validation runs docker or envoy reaches things; keep default.
# If you later validate only locally, you can tighten further.
